# deployment/valetudo.conf
#!upstart
description "Valetudo"

start on filesystem and net-device-up IFACE=lo
stop on runlevel [06]

oom score 1000

script
	# give other software some time to initialize
	echo -n "Waiting for 15 sec.."
	sleep 15
	echo " done."
	# check data partition to be mounted before doing anything
	while ! /bin/mountpoint -q /mnt/data
	do
		echo "Data mountpoint isn't ready, can't run yet. Retrying in 5 sec..."
		sleep 5
	done
	# disable core dumps on this system, we're in production
	echo '|/bin/false' > /proc/sys/kernel/core_pattern
	# finally run valetudo
	exec /usr/local/bin/valetudo
end script

respawn
respawn limit 10 90
limit as 209715200 209715200
